# container/cron.d/logRotater
*/15 * * * * phvalheim /opt/stateless/engine/tools/logRotater.sh >> /opt/stateful/logs/logRotater.log
