# container/cron.d/tsStoreSync
0 */12 * * * phvalheim /opt/stateless/engine/tools/tsStoreSync.sh
