# container/cron.d/tsSyncLocalParse
0 */12 * * * phvalheim /opt/stateless/engine/tools/tsSyncLocalParseMultithreaded.sh >> /opt/stateful/logs/tsSync.log
