# container/cron.d/tsSyncRemoteParse
15 * * * * phvalheim /opt/stateless/engine/tools/tsSyncRemoteParse.sh >> /opt/stateful/logs/phvalheim.log
