# container/cron.d/utilizationMonitor
* * * * * root /opt/stateless/engine/tools/utilizationMonitor >> /opt/stateful/logs/phvalheim.log
