# container/cron.d/world_backup
*/30 * * * * phvalheim /opt/stateless/engine/tools/world_backup
