# container/cron.d/worldBackup
0 * * * * phvalheim /opt/stateless/engine/tools/worldBackup >> /opt/stateful/logs/worldBackups.log
