# container/cron.d/worldBackupRetention
45 * * * * phvalheim /opt/stateless/engine/tools/worldBackupRetention >> /opt/stateful/logs/worldBackups.log 2>&1
