# container/cron.d/worldBackup
*/30 * * * * phvalheim /opt/stateless/engine/tools/worldBackup >> /opt/stateful/logs/phvalheim.log
