# container/cron.d/worldActivityMonitor
*/5 * * * * phvalheim /opt/stateless/engine/tools/worldActivityMonitor >> /opt/stateful/logs/phvalheim.log 2>&1
