# container/cron.d/worldBackupCompress
30 * * * * phvalheim /opt/stateless/engine/tools/worldBackupCompress >> /opt/stateful/logs/worldBackups.log 2>&1
